# 6min.target
[Unit]
Description=6min ticker target
StopWhenUnneeded=yes

[Install]
Also=6min.timer
WantedBy=multi-user.target
